# units/cpu-eater-large.service
[Unit]
Description=Eat CPU

[Service]
Restart=always
CPUShares=1500
ExecStart=/bin/bash -c 'while true; do echo foo; done > /dev/null'
